# com.cast.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  CAST Background Daemon — macOS LaunchAgent
  Install: cp com.cast.daemon.plist ~/Library/LaunchAgents/
  Load:    launchctl load ~/Library/LaunchAgents/com.cast.daemon.plist
  Stop:    launchctl unload ~/Library/LaunchAgents/com.cast.daemon.plist
  Logs:    tail -f ~/.claude/logs/castd.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cast.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/edkubiak/.claude/scripts/castd.sh</string>
  </array>

  <!-- Start immediately on login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart automatically on crash; min 10s between restarts -->
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Log stdout + stderr to shared log file -->
  <key>StandardOutPath</key>
  <string>/Users/edkubiak/.claude/logs/castd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/edkubiak/.claude/logs/castd.log</string>

  <!-- Inherit a login-like environment so PATH includes Homebrew + pyenv -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/edkubiak</string>
  </dict>
</dict>
</plist>
